# deploy/supervisor/gym2-reverb.conf
; WebSocket server for real-time dashboard/notification broadcasting.
; Install to /etc/supervisor/conf.d/, then:
;   supervisorctl reread && supervisorctl update && supervisorctl start gym2-reverb
;
; Replace /var/www/gym2 and the user below to match the deployment host.
; Reverb's own bind host/port/hostname come from .env (REVERB_HOST/PORT) —
; put a reverse proxy (Nginx) in front for TLS rather than exposing this
; port directly.

[program:gym2-reverb]
directory=/var/www/gym2
command=php artisan reverb:start
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/gym2/storage/logs/reverb.log
