# deploy/supervisor/gym2-worker.conf
; Queue worker for the database queue (notifications, broadcast listeners —
; see app/Listeners). Install to /etc/supervisor/conf.d/, then:
;   supervisorctl reread && supervisorctl update && supervisorctl start gym2-worker:*
;
; Replace /var/www/gym2 and the user below to match the deployment host.

[program:gym2-worker]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www/gym2
command=php artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/gym2/storage/logs/worker.log
stopwaitsecs=3600
